# deploy/systemd/aevia-node.service
[Unit]
Description=Aevia Node (Provider / Relay — libp2p + DHT + WHIP + optional HTTPS)
Documentation=https://github.com/Leeaandrob/aevia
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aevia
Group=aevia
WorkingDirectory=/var/lib/aevia-node
# EnvironmentFile drives the whole config — AEVIA_MODE, AEVIA_LISTEN,
# AEVIA_HTTP_ADDR, AEVIA_BOOTSTRAP, AEVIA_RELAY_PEERS, AEVIA_TLS_DOMAIN,
# AEVIA_CF_API_TOKEN, etc. See deploy/systemd/env.example.
EnvironmentFile=/etc/aevia-node/env
ExecStart=/usr/local/bin/aevia-node run --data-dir=/var/lib/aevia-node

Restart=always
RestartSec=5s
TimeoutStopSec=30s

# Allow binding to privileged ports (443, 80) without running as root —
# needed by provider nodes that serve HTTPS. Relay nodes that only bind
# 4001/8080 can drop this cap with no effect.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/aevia-node
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

# Resource caps (provider-node with WHIP ingest + CMAF segmenter + TLS
# needs more headroom than a relay; keep MemoryMax generous).
LimitNOFILE=65536
MemoryHigh=1G
MemoryMax=2G

[Install]
WantedBy=multi-user.target
